# cookbooks/polartour_workspace/templates/default/polartour.conf.erb
description "polartour"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on [!12345]

console log

env HOME=<%= node[:project_workspace][:dir] %>
setuid <%= node[:project_workspace][:user] %>
setgid <%= node[:project_workspace][:group] %>
chdir <%= node[:project_workspace][:dir] %>/polartour
exec coffee index.coffee
